# deploy/supervisor/app.conf
# =====================================================
# Supervisor Configuration for MentorHub
# Optimized for Graceful Shutdown
# Works on: Render, Railway, Fly.io, VPS, Amvera, Heroku
# =====================================================

[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/var/run/supervisord.pid
user=root
loglevel=info

# =====================================================
# Backend (FastAPI + Uvicorn)
# =====================================================
[program:backend]
# Uvicorn с настройками graceful shutdown
# timeout-graceful-shutdown: время на завершение активных запросов
# Используем $PORT из переменной окружения (Render и другие платформы)
command=uvicorn app.main:app --host 0.0.0.0 --port %(ENV_PORT)s --workers 2 --timeout-graceful-shutdown 30
directory=/app/backend
autostart=true
autorestart=true

# КРИТИЧНО: правильная остановка группы процессов
# stopasgroup=true - отправляет SIGTERM всей группе процессов (включая дочерние workers)
# killasgroup=true - отправляет SIGKILL всей группе если stopwaitsecs истёк
stopasgroup=true
killasgroup=true
stopwaitsecs=35

# Логи в Docker stdout/stderr для просмотра через docker logs
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Приоритет запуска (backend стартует первым)
priority=10

# Environment variables - PORT передаётся извне (Render устанавливает $PORT)
environment=PORT="%(ENV_PORT)s",PYTHONUNBUFFERED="1"

# =====================================================
# Frontend (Next.js Standalone)
# =====================================================
[program:frontend]
command=node server.js
directory=/app/frontend
environment=PORT="3000",NODE_ENV="production"
autostart=true
autorestart=true

# Graceful shutdown для Node.js
stopasgroup=true
killasgroup=true
stopwaitsecs=20

# Логи в Docker stdout/stderr
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Приоритет запуска (frontend стартует после backend)
priority=20
